# systemd-services/cyclic-sysinfo.service
[Unit]
Description=cyclic-sysinfo

[Service]
ExecStart=/usr/local/bin/cyclic-sysinfo-exec -c 10
Type=simple

[Install]
WantedBy=multi-user.target
