# systemd-services/epoll-psi.service
[Unit]
Description=epoll-psi

[Service]
ExecStart=/usr/local/bin/epoll-psi-exec -c 10
Type=simple

[Install]
WantedBy=multi-user.target
