# systemd-services/inotify-coredump.service
[Unit]
Description=inotify-coredump

[Service]
ExecStart=/usr/local/bin/inotify-coredump-exec
Type=simple

[Install]
WantedBy=multi-user.target
