# systemd-services/posix-mq-to-mqtt-bridge.service
[Unit]
Description=posix-mq-to-mqtt-bridge

[Service]
ExecStart=/usr/local/bin/bridge-exec
Type=simple

[Install]
WantedBy=multi-user.target
